# files/etc/uci-defaults/open-wireless
sed -i "/option disabled '1'/d" /etc/config/wireless
echo "All done!"
