# files/etc/uci-defaults/set_usbwan-wireless_enadbled
uci set wireless.radio0.disabled='0'
uci set wireless.radio1.disabled='0'
uci set wireless.radio1.channel='64'
uci set wireless.radio1.country='RU'
uci commit wireless
uci set network.usbwan=interface
uci set network.usbwan.proto='dhcp'
uci set network.usbwan.device='eth1'
uci set network.usbwan6=interface
uci set network.usbwan6.proto='dhcpv6'
uci set network.usbwan6.device='eth1'
uci commit network                  
echo "All done!"
